# scripts/backup.timer

[Unit]
Description= Runs backup_root_fs.sh on Fridays at 21:00

[Timer]
OnCalendar=Fri 21:00
Unit= backup.service

[Install]
WantedBy= basic.target